# standx-account1.service
[Unit]
Description=StandX Account 1
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/standx
ExecStart=/root/standx/.venv/bin/python -u /root/standx/market_maker.py --config /root/standx/.env.account1 --log-prefix account1
ExecStop=/bin/bash -c 'kill -TERM $(pgrep -f "market_maker.py --config.*account1")'

Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
