# standx-market-maker.service
[Unit]
Description=StandX Market Maker
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/standx
ExecStart=/root/standx/run.sh
ExecStop=/root/standx/stop.sh

# --- 守护进程自动重启配置 ---
Restart=always
RestartSec=15
StartLimitIntervalSec=300
StartLimitBurst=5
TimeoutStopSec=45
KillMode=mixed
KillSignal=SIGTERM

# --- 资源限制开始 ---
# 限制 CPU 使用率为单核的 35%
CPUQuota=20%
# 内存使用达到 500M 时，系统会尝试回收内存页，进程仍可运行但性能可能受限
MemoryHigh=500M
# 内存硬上限 800M，超过此值系统将直接 Kill 进程（保护节点不重启）
MemoryMax=800M
# 限制该服务及其子进程的总线程/进程数
TasksMax=50
# 禁止使用 Swap 交换分区，确保内存使用的确定性（可选，推荐）
MemorySwapMax=0
# --- 资源限制结束 ---
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target